# linux_files/dotfiles/.config/autostart/terminator.desktop
[Desktop Entry]
Exec=terminator
Name=terminator
Type=Application
Version=1.0
